# infra/ai.openclaw.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.update</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/ainchorsangiefpl/.openclaw/workspace/infra/openclaw-update-wrapper.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/openclaw-update-2026-07-01-v3.plist.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/openclaw-update-2026-07-01-v3.plist.err</string>
</dict>
</plist>
